# systemd/tricky-mesh-ai.service
[Unit]
Description=Tricky Mesh AI bridge (Remote-Terminal webhook -> llama.cpp)
# Start after network + llama.cpp model server. llama-server is Type=simple
# and doesn't signal readiness, so we only get "start-after" semantics, not
# "wait-until-ready" — during the ~30s Qwen 3.6-35B model load, early DMs
# may fail their LLM call and land in the dead-letter log. Subsequent DMs
# are fine once the model is up. Good enough for boot-time ordering without
# extra machinery.
#
# RT (remoteterm.service) is system-scope and always starts before any
# user-scope service via the multi-user.target → default.target chain, so
# no explicit ordering needed for RT here.
After=network-online.target llama-server.service
Wants=network-online.target llama-server.service

[Service]
Type=simple
ExecStart=%h/projects/mesh/tricky-mesh-ai-api/.venv/bin/tricky-mesh-ai
Restart=on-failure
RestartSec=5
# Don't restart if the daemon exits because of a config error (exit 2).
RestartPreventExitStatus=2
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tricky-mesh-ai

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/tricky-mesh-ai-api
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=default.target
